# Raspberry/levara.service.d/rerank.conf
# Drop-in: /etc/systemd/system/levara.service.d/rerank.conf (живёт на Pi, проверено 2026-05-27).
# Подключает reranker (mmini-L12-int8) к Levara.
[Service]
Environment=RERANK_ENDPOINT=http://localhost:9100/rerank
Environment=RERANK_MODEL=mmini-L12-int8
Environment=RERANK_BUDGET_MS=5000
Environment=RERANK_TIMEOUT_MS=5000
