# deploy/bench/embed-bench.service
[Unit]
Description=Embed-bench sidecar (Levara load-profile model swap)
After=network-online.target

[Service]
Type=simple
User=stek0v
WorkingDirectory=/home/stek0v/embed-bench/scripts/load-profiles
Environment=PYTHONUNBUFFERED=1
Environment=HF_HOME=/home/stek0v/embed-bench/hf-cache
ExecStart=/home/stek0v/embed-bench/venv/bin/uvicorn embed_bench.server:app --host 127.0.0.1 --port 9101 --workers 1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
