# Raspberry/levara.service
# Реальный деплой на Pi (10.23.0.53), node-id=pi1, проверено 2026-05-27.
# Установка: /etc/systemd/system/levara.service (+ drop-in levara.service.d/rerank.conf).
[Unit]
Description=Levara MCP Knowledge Engine
After=network-online.target ollama.service
Wants=network-online.target ollama.service

[Service]
Type=simple
User=stek0v
WorkingDirectory=/home/stek0v/levara
Environment=DB_PROVIDER=sqlite
Environment=DB_PATH=/home/stek0v/levara/data/levara.db
Environment=EMBEDDING_ENDPOINT=http://localhost:11434/v1/embeddings
Environment=EMBEDDING_MODEL=nomic-embed-text-v2-moe
Environment=LLM_PROVIDER=openai
Environment=LLM_ENDPOINT=http://localhost:11434/v1
Environment=LLM_MODEL=qwen3:0.6b
Environment=VISION_ENDPOINT=http://10.23.0.101:8081/api/v1/ocr
Environment=LOG_LEVEL=INFO
ExecStartPre=/bin/sleep 5
ExecStart=/home/stek0v/levara/levara -standalone=true -dim=768 -shards=1 -port=8090 -grpc-port=0 -data-dir=/home/stek0v/levara/data -node-id=pi1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
